# service/webui/run
#!/bin/sh
# Run the Flask web UI on port 5000
exec python3 /home/dockeruser/app.py